# profile.d/dell4khack.sh
export LD_PRELOAD=/usr/local/\$LIB/dell4khack.so
