# profile.d/hax11.sh
export LD_PRELOAD=/usr/local/\$LIB/hax11.so
