# profile.d/mst4khack.sh
export LD_PRELOAD=/usr/local/\$LIB/mst4khack.so
